# puls8_dbf_data_prep/schedule/custom_est.cron
CRON_TZ=America/New_York
# Example: Every Saturday at 8 AM EST
0 8 * * 6 docker run --rm -v /home/ec2-user/jobs/logs:/logs puls8_dbf_data_prep:latest
